# provision/systemd/update-docs-ja.timer
[Unit]
Description=Update docs.ruby-lang.org/ja/

[Timer]
OnCalendar=*-*-* 20:15
Persistent=true
Unit=bc-setup-all.service

[Install]
WantedBy=timers.target
